# config/systemd/rango.service.erb
[Unit]
Description = <%= fetch(:application) %> (<%= fetch(:stage) %>) server
After=syslog.target network.target

[Service]
WatchdogSec=10
Type=simple
Environment=PWD=<%= current_path %>
EnvironmentFile=<%= current_path %>/.env
WorkingDirectory=<%= current_path %>
ExecStart=<%= current_path %>/rango

KillMode=mixed
Restart=always
RestartSec=1
SyslogIdentifier=rango

[Install]
# Use default.target to autoboot user service
# https://superuser.com/questions/1025091/start-a-systemd-user-service-at-boot
# WantedBy = multi-user.target
WantedBy=default.target
